# deployment/systemd/gps-receivers-scheduler.service
[Unit]
# This is a USER-LEVEL systemd unit, installed at
#   ~gpsops/.config/systemd/user/gps-receivers-scheduler.service
# Operator manages it as gpsops (no sudo):
#   systemctl --user restart gps-receivers-scheduler
#   journalctl --user-unit gps-receivers-scheduler -f
# Linger must be enabled so the unit runs without an active gpsops session:
#   loginctl enable-linger gpsops
# install.sh handles all of the above.

Description=GPS Receivers Pipeline Scheduler - Automated GNSS Data Processing
Documentation=https://github.com/bennigo/receivers
# Note: user units cannot Require= system units (postgresql, network-online).
# After= still gives ordering on boot when linger triggers user-manager startup.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/bgo/git/receivers

Environment="LANG=en_US.UTF-8"
Environment="GPS_CONFIG_PATH=/home/gpsops/.config/gpsconfig"

# Pipeline scheduler schedules ALL stations from stations.cfg.
# Worker count is controlled via `scheduler.max_workers` in scheduler.yaml —
# no --max-workers CLI flag here, so a single knob governs it.
# Use --stations flag on manual CLI runs to test subsets.
ExecStart=/home/bgo/git/receivers/venv/bin/receivers scheduler start

# Graceful shutdown - allow active pipelines to complete
ExecStop=/home/bgo/git/receivers/venv/bin/receivers scheduler stop
TimeoutStopSec=120

# Restart policy
Restart=always
RestartSec=30s
StartLimitInterval=10min
StartLimitBurst=5

# Resource limits — for user units these only enforce when cgroup delegation
# is enabled on user@.service. install.sh installs the
# /etc/systemd/system/user@.service.d/delegate.conf drop-in for that.
# If you ever copy this unit to a host without that drop-in, MemoryMax and
# CPUQuota will be silently ignored.
LimitNOFILE=65536
MemoryMax=4G
MemoryHigh=3G
CPUQuota=400%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gps-scheduler

# Security
# NOTE: filesystem-sandbox directives (PrivateTmp, ProtectSystem, ProtectHome,
# ReadWritePaths) are intentionally NOT used here. For a `systemctl --user`
# service they force a USER NAMESPACE in which root-owned files appear as
# `nobody` (overflow uid 65534). OpenSSH then rejects /etc/ssh/ssh_config.d/*.conf
# with "Bad owner or permissions" (its st_uid != 0 check), so rsync-over-ssh to
# the long-term archive — push-on-download AND the :45 archive-sync — fails with
# rc=255. The scheduler already runs unprivileged as gpsops; NoNewPrivileges and
# the resource limits above are kept (they need no mount namespace). See the
# 2026-06-21 archive-cutover incident.
NoNewPrivileges=true

# Watchdog — DISABLED.
# WatchdogSec only works when Type=notify AND the application sends
# sd_notify WATCHDOG=1 keep-alives. The scheduler is Type=simple and
# doesn't send keep-alives, so enabling this puts systemd in a kill loop
# (every 5 min SIGKILL → restart) that prevents any job from ever firing.
# Re-enable only after integrating sd_notify in the scheduler main loop
# (e.g. via systemd-python or `sdnotify` package) AND flipping Type=notify.
#WatchdogSec=300s
#NotifyAccess=main

[Install]
# default.target is the user-unit equivalent of multi-user.target.
WantedBy=default.target
