# deployment/systemd/gps-receivers-scheduler.service
[Unit]
Description=GPS Receivers Pipeline Scheduler - Automated GNSS Data Processing
Documentation=https://github.com/bennigo/receivers
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=gpsops
Group=gpsops
WorkingDirectory=/opt/receivers

# Environment setup
Environment="GPS_CONFIG_PATH=/etc/gpsconfig"
Environment="GPS_CACHE_DIR=/var/cache/gps_receivers"
Environment="LANG=en_US.UTF-8"

# Command to run - Pipeline scheduler schedules ALL stations from stations.cfg
# Use --stations flag only for testing subsets
ExecStart=/opt/receivers/venv/bin/receivers scheduler start \
    --max-workers 5

# Graceful shutdown - allow active pipelines to complete
ExecStop=/opt/receivers/venv/bin/receivers scheduler stop
TimeoutStopSec=120

# Restart policy
Restart=always
RestartSec=30s
StartLimitInterval=10min
StartLimitBurst=5

# Resource limits
# Network pool: 10 workers (I/O-bound, concurrent OK)
# CPU pool: 4 workers (memory-limited for RINEX conversion)
LimitNOFILE=65536
MemoryMax=4G
MemoryHigh=3G
CPUQuota=400%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gps-scheduler

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/var/cache/gps_receivers /mnt/gpsdata /tmp

# Watchdog (restart if no keep-alive for 5 minutes)
WatchdogSec=300s
NotifyAccess=main

[Install]
WantedBy=multi-user.target
